# iso-build/config/includes.chroot/etc/skel/.config/autostart/circleos-shell.desktop
[Desktop Entry]
Type=Application
Name=CircleOS Shell
Exec=circleos-shell
Hidden=false
NoDisplay=false
X-GNOME-Autostart-enabled=true
